# deploy/systemd/galactic-sovereign-coop.service
[Unit]
Description=Galactic Sovereign persistent multiplayer universe
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
Type=simple
User=galactic-sovereign
Group=galactic-sovereign
WorkingDirectory=/opt/galactic-sovereign/current
ExecStart=/usr/bin/node /opt/galactic-sovereign/current/server/coop-host.mjs
Environment=NODE_ENV=production
Environment=GS_COOP_HOST=127.0.0.1
Environment=GS_COOP_PORT=9090
Environment=GS_COOP_DATA_DIR=/var/lib/galactic-sovereign/multiplayer
Environment=GS_COOP_MAX_PLAYERS=5
Environment=GS_COOP_AUTOSAVE_MS=60000
LoadCredential=gateway-secret:/etc/galactic-sovereign/credentials/gateway-secret
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/var/lib/galactic-sovereign/multiplayer
UMask=0077
MemoryHigh=450M
MemoryMax=650M
TasksMax=128

[Install]
WantedBy=multi-user.target
