# deploy/systemd/galactic-sovereign-backup.timer
[Unit]
Description=Hourly Galactic Sovereign backups

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true
Unit=galactic-sovereign-backup.service

[Install]
WantedBy=timers.target
