# deploy/systemd/galactic-sovereign-gateway.service
[Unit]
Description=Galactic Sovereign authenticated gateway
After=network-online.target galactic-sovereign-coop.service
Wants=network-online.target galactic-sovereign-coop.service
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
Type=simple
User=galactic-sovereign
Group=galactic-sovereign
WorkingDirectory=/opt/galactic-sovereign/current
ExecStart=/usr/bin/node /opt/galactic-sovereign/current/server/app-host.mjs
Environment=NODE_ENV=production
Environment=GS_APP_HOST=127.0.0.1
Environment=GS_APP_PORT=8080
Environment=GS_PUBLIC_ORIGIN=https://play.galacticsovereign.xyz
Environment=GS_DATA_DIR=/var/lib/galactic-sovereign/accounts
Environment=GS_DIST_DIR=/opt/galactic-sovereign/current/dist
Environment=GS_COOP_INTERNAL_URL=ws://127.0.0.1:9090
LoadCredential=gateway-secret:/etc/galactic-sovereign/credentials/gateway-secret
LoadCredential=session-pepper:/etc/galactic-sovereign/credentials/session-pepper
Restart=always
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/var/lib/galactic-sovereign/accounts
UMask=0077
MemoryHigh=350M
MemoryMax=500M
TasksMax=128

[Install]
WantedBy=multi-user.target
